# .claude/skills/launchd/com.trade-analysis.skill-generation-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  SETUP: Replace $HOME and $PROJECT_DIR before installing.
    sed "s|\$HOME|$HOME|g; s|\$PROJECT_DIR|$(pwd)|g" \
      launchd/com.trade-analysis.skill-generation-daily.plist > ~/Library/LaunchAgents/com.trade-analysis.skill-generation-daily.plist
    launchctl load ~/Library/LaunchAgents/com.trade-analysis.skill-generation-daily.plist
-->
<dict>
	<key>Label</key>
	<string>com.trade-analysis.skill-generation-daily</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>$PROJECT_DIR/scripts/run_skill_generation.sh</string>
		<string>--mode</string>
		<string>daily</string>
	</array>

	<key>WorkingDirectory</key>
	<string>$PROJECT_DIR</string>

	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>7</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/opt/homebrew/sbin:$HOME/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>HOME</key>
		<string>$HOME</string>
	</dict>

	<key>StandardOutPath</key>
	<string>$PROJECT_DIR/logs/launchd_skill_generation_daily.log</string>

	<key>StandardErrorPath</key>
	<string>$PROJECT_DIR/logs/launchd_skill_generation_daily_error.log</string>

	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
